# .config/systemd/user/media-sync.timer
[Unit]
Description=Media Sync timer

[Timer]
OnCalendar=*-*-* *:*:00
Unit=media-sync.service

[Install]
WantedBy=multi-user.target
